# dotfiles/profile.d/adium.env
defaults write com.adiumx.Adiumx AMXMPPShowAdvanced -bool YES
