# dotfiles/profile.d/adium.env
if [ "${PLATFORM}" = "Darwin" ]; then
  defaults write com.adiumx.Adiumx AMXMPPShowAdvanced -bool YES
fi
