# dotfiles/profile.d/environment
source ${HOME}/.profile.d/global.functions
source ${HOME}/.profile.d/global.env
source ${HOME}/.profile.d/global.aliases

source ${HOME}/.profile.d/work.functions

source ${HOME}/.profile.d/all.aliases
source ${HOME}/.profile.d/personal.aliases
source ${HOME}/.profile.d/oss.aliases
source ${HOME}/.profile.d/haskell.aliases
source ${HOME}/.profile.d/ruby.aliases
source ${HOME}/.profile.d/scala.aliases
source ${HOME}/.profile.d/work.aliases

source ${HOME}/.profile.d/clojure.env
source ${HOME}/.profile.d/ec2.env
source ${HOME}/.profile.d/erlang.env
source ${HOME}/.profile.d/haskell.env
source ${HOME}/.profile.d/github.env
source ${HOME}/.profile.d/java.env
source ${HOME}/.profile.d/javascript.env
source ${HOME}/.profile.d/jruby.env
source ${HOME}/.profile.d/node.env
source ${HOME}/.profile.d/opa.env
source ${HOME}/.profile.d/python.env
source ${HOME}/.profile.d/ruby.env
source ${HOME}/.profile.d/scala.env
source ${HOME}/.profile.d/ssh.env
source ${HOME}/.profile.d/personal.env
source ${HOME}/.profile.d/work.env
source ${HOME}/.profile.d/adium.env

if [ "${PLATFORM}" == "Darwin" ]; then
  source ${HOME}/.profile.d/osx.env
  source ${HOME}/.profile.d/osx.functions
  source ${HOME}/.profile.d/osx.aliases
fi;

# User specific environment and startup programs

export PATH EDITOR GAEPATH PYTHONPATH RAILS_ENV MOZ_DISABLE_PANGO RUBYLIB
export ERL_LIBS AUTOFEATURE RSPEC SSH_AUTH_SOCK SDL_AUDIODRIVER LANG LC_ALL
export PS_PERSONALITY PROMPT PS2 RPROMPT WINDOW_MANAGER DIALYZER_PLT
export PLATFORM
