# dotfiles/profile.d/environment
function sourceProfileInit() {
  if [ "$#" -ne 1 ]; then
    >&2 echo "Error: Expected one argument, received $# for ${FUNCNAME}"
    >&2 echo "Usage: ${FUNCNAME} FILE_INSIDE_PROFILE_DIR"
    return 1
  fi

  local -r profile_path="${HOME}/.profile.d/${1}"

  if [ -f "${profile_path}" ]; then
    source "${profile_path}"
  fi
}

sourceProfileInit "global.functions"
sourceProfileInit "global.env"
sourceProfileInit "global.aliases"

sourceProfileInit "work.functions"
sourceProfileInit "personal.functions"
sourceProfileInit "oss.functions"
sourceProfileInit "haskell.functions"
sourceProfileInit "idris.functions"
sourceProfileInit "nix.functions"
sourceProfileInit "erlang.functions"
sourceProfileInit "scala.functions"
sourceProfileInit "ruby.functions"
sourceProfileInit "python.functions"
sourceProfileInit "github.functions"
sourceProfileInit "ssh.functions"
sourceProfileInit "gerrit.functions"
sourceProfileInit "github.functions"
sourceProfileInit "lookout.functions"
sourceProfileInit "crypto.functions"

sourceProfileInit "nix.env"
sourceProfileInit "clojure.env"
sourceProfileInit "ec2.env"
sourceProfileInit "erlang.env"
sourceProfileInit "haskell.env"
sourceProfileInit "idris.env"
sourceProfileInit "github.env"
sourceProfileInit "java.env"
sourceProfileInit "javascript.env"
sourceProfileInit "jruby.env"
sourceProfileInit "node.env"
sourceProfileInit "opa.env"
sourceProfileInit "python.env"
sourceProfileInit "ruby.env"
sourceProfileInit "scala.env"
sourceProfileInit "ssh.env"
sourceProfileInit "personal.env"
sourceProfileInit "work.env"
sourceProfileInit "adium.env"
sourceProfileInit "gerrit.env"
sourceProfileInit "lookout.env"
sourceProfileInit "xmonad.env"
sourceProfileInit "gradle.env"

sourceProfileInit "all.aliases"
sourceProfileInit "personal.aliases"
sourceProfileInit "oss.aliases"
sourceProfileInit "haskell.aliases"
sourceProfileInit "idris.aliases"
sourceProfileInit "ruby.aliases"
sourceProfileInit "scala.aliases"
sourceProfileInit "work.aliases"
sourceProfileInit "gerrit.aliases"
sourceProfileInit "lookout.aliases"
sourceProfileInit "linux.aliases"

if [ "${PLATFORM}" = "Darwin" ]; then
  sourceProfileInit "osx.env"
  sourceProfileInit "osx.functions"
  sourceProfileInit "osx.aliases"
fi;

if [ "${PLATFORM}" = "Linux" ]; then
  sourceProfileInit "linux.env"
fi

export PATH EDITOR GAEPATH PYTHONPATH RAILS_ENV MOZ_DISABLE_PANGO RUBYLIB
export ERL_LIBS AUTOFEATURE RSPEC SSH_AUTH_SOCK SDL_AUDIODRIVER LANG LC_ALL
export PS_PERSONALITY PROMPT PS2 RPROMPT WINDOW_MANAGER DIALYZER_PLT
export PLATFORM
