# dotfiles/profile.d/erlang.env
# Erlang variables:
DIALYZER_PLT=${HOME}/.dialyzer

export DIALYZER_PLT
