# dotfiles/profile.d/gerrit.functions

function gerrit() {
  ssh -p "${GERRIT_PORT}" "${GERRIT_HOST}" gerrit $*
}
