# dotfiles/profile.d/gerrit.env
# Define GERRIT_HOST and GERRIT_PORT here if not sensitive
# Otherwise in hidden *.env file :)
