# dotfiles/profile.d/gerrit.aliases
alias gerrit="ssh -p 29418 gerrit.flexilis.local gerrit"
