# dotfiles/profile.d/haskell.env
PATH=${HOME}/.cabal/bin:${PATH}
GHC_OPTS=development

export PATH GHC_OPTS
