# dotfiles/profile.d/haskell.env
PATH=${HOME}/.cabal/bin:${PATH}:.cabal-sandbox/bin
GHC_OPTS=development

export PATH GHC_OPTS
