# dotfiles/profile.d/haskell.env
PATH=${PATH}:${HOME}/.cabal/bin
GHC_OPTS=development

export PATH GHC_OPTS
