# dotfiles/profile.d/haskell.env
PATH=${PATH}:${HOME}/.cabal/bin

export PATH
