# dotfiles/profile.d/haskell.aliases
#alias cabal=cabal-dev
alias ydev="yesod devel"
