# dotfiles/profile.d/javascript.env
PATH=${PATH}:${HOME}/.npmjs/bin

export PATH
