# dotfiles/profile.d/jruby.env
JRUBY_PATH=${HOME}/opt/jruby

[ -d ${JRUBY_PATH} ] && PATH=${PATH}:${JRUBY_PATH}/bin

export PATH
