# dotfiles/profile.d/jruby.env
JRUBY_PATH=${HOME}/opt/jruby
PATH=${PATH}:${JRUBY_PATH}/bin

export PATH
