# dotfiles/profile.d/linux.env
SCREENSHOTS_DIR="${HOME}/.screenshots"

mkdir -p "${SCREENSHOTS_DIR}"

export SCREENSHOTS_DIR
