# dotfiles/profile.d/nix.env
NIX_BASE="${HOME}/.nix-profile"
NIX_PKGS="${NIX_BASE}/nixpkgs"
NIX_REMOTE=daemon
#unset NIX_REMOTE
NIX_PATH="${NIX_BASE}/bin:${NIX_BASE}/usr/bin:nixpkgs=${NIX_PKGS}"

export NIX_BASE NIX_PATH PATH

