# dotfiles/profile.d/node.env
NPM_DIR=${HOME}/.npmjs

PATH=${PATH}:${NPM_DIR}/bin
export NPM_DIR PATH
