# dotfiles/profile.d/node.env
NPM_DIR=${HOME}/.npmjs
NPM_BINDIR=${NPM_DIR}/bin

[ -d ${NPM_BINDIR} ] && PATH=${PATH}:${NPM_BINDIR}

export NPM_DIR PATH
