# dotfiles/profile.d/opa.env
PATH=/opt/mlstate/bin:${PATH}

export PATH
