# dotfiles/profile.d/opa.env
OPAPATH=${HOME}/opt/mlstate/bin

[ -d ${OPAPATH} ] && PATH=${OPAPATH}:${PATH}

export PATH
