# dotfiles/profile.d/personal.env
CODE_DIR=${HOME}/src
DEFAULT_BROWSER=firefox

export CODE_DIR DEFAULT_BROWSER

