# dotfiles/profile.d/python.env
GAEPATH=${HOME}/opt/gae
DJANGOPATH=${GAEPATH}/lib/django

[ -d ${GAEPATH} ] && PATH=${GAEPATH}:${PATH} && \
  PYTHONPATH=${GAEPATH}:${PYTHONPATH}
[ -d ${DJANGOPATH} ] && PATH=${DJANGOPATH}/django/bin:${PATH} && \
  PYTHONPATH=${DJANGOPATH}:${PYTHONPATH}

VIRTUALENVWRAPPER_PYTHON="$(readlink $(which python))"
PYTHONPATH=${PYTHONPATH}:.:tests
WORKON_HOME=${HOME}/.pyenvs
mkdir -p ${WORKON_HOME}

export PYTHONPATH PATH WORKON_HOME VIRTUALENVWRAPPER_PYTHON
