# dotfiles/profile.d/python.env
GAEPATH=${HOME}/opt/gae
DJANGOPATH=${GAEPATH}/lib/django
PYTHONPATH=${GAEPATH}:${DJANGOPATH}:${PYTHONPATH}:.:tests
PATH=${GAEPATH}:${DJANGOPATH}/django/bin

export DJANGOPATH PYTHONPATH PATH
