# dotfiles/profile.d/ruby.env
RUBYLIB=lib
AUTOFEATURE=true
RSPEC=true
JRUBY_OPTS=--1.9
#RAILS_ENV=development

PATH=${HOME}/.rbenv/shims:${PATH}

export PATH RUBYLIB AUTOFEATURE RSPEC JRUBY_OPTS
