# dotfiles/profile.d/ruby.env
RUBYLIB=lib
AUTOFEATURE=true
RSPEC=true
#RAILS_ENV=development

PATH=${HOME}/.rvm/rubies/default/bin:${PATH}

export PATH RUBYLIB AUTOFEATURE RSPEC
