# dotfiles/profile.d/ruby.aliases
# Bundle prefix aliases
alias brake="bundle exec rake"

# JRuby aliases
alias jgem="jruby -S gem"
alias jrails="jruby -S rails"
alias jrake="jruby -S rake"

