# dotfiles/profile.d/ruby.aliases
# Bundle prefix aliases
alias bx="bundle exec"
alias brake="bundle exec rake"
alias bguard="bundle exec guard"

# JRuby aliases
alias jgem="jruby -S gem"
alias jrails="jruby -S rails"
alias jrake="jruby -S rake"

