# dotfiles/profile.d/scala.env
# Scala paths
SCALA_ROOT=/usr/share/scala
SCALA_LIB=${SCALA_ROOT}/lib

export SCALA_ROOT SCALA_LIB
