# dotfiles/profile.d/ssh.env
SSHAGENT=/usr/bin/ssh-agent
SSH_AUTH_SOCK=${HOME}/.ssh/agent.sock
SSHAGENTARGS="-a ${SSH_AUTH_SOCK}"

export SSHAGENT SSHAGENTARGS SSH_AUTH_SOCK
