# dotfiles/profile.d/work.functions

function lookout_vpn_connect() {
  sudo vpnc "${HOME}/.vpn.d/lookoutops.vpnc"
}
