# dotfiles/profile.d/xmonad.env
if [ ! -z "${XMONAD_GHC}" ]; then
  XMONAD_SCREENSHOT_DIR="${HOME}/.screenshots"
  mkdir -p "${XMONAD_SCREENSHOT_DIR}"

  export XMONAD_SCREENSHOT_DIR
fi
